# docker/crontab
* * * * * root cd /var/www/html && bash crontab.sh >> /dev/null 2>&1
